# docker/supervisor/home_assistant.ini
[program:web_server]
command=/app/bin/server
autostart=true
stdout_logfile=/app/log/web-server.stdout.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
stderr_logfile=/app/log/web-server.stderr.log
stderr_logfile_maxbytes=50MB
stderr_logfile_backups=10
priority=4
